# config/selenium.service
[Unit]
Description=Manage a selenium grid hub
Documentation=https://github.com/thejsninja/vagrant-selenium-hub

[Service]
ExecStart=/usr/bin/java -jar /home/vagrant/selenium-standalone-server.jar -role hub
Type=simple

[Install]
WantedBy=multi-user.target
